# etc/init/check-rw-vpd.conf
# Copyright 2016 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description "Check RW_VPD partition and report metrics on read status"
author      "chromium-os-dev@chromium.org"

start on started system-services

exec check_rw_vpd
