# etc/init/cups-post-upstart-socket-bridge.conf
# Copyright 2016 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Prepare for CUPS, after upstart-socket-bridge"

start on started upstart-socket-bridge

env CUPS_SOCKET=/run/cups/cups.sock
env lp_group=lp

# HACK: Set permissions, since upstart-socket-bridge will create according to
# umask (i.e., 755), which means non-root users can't talk to cups...
script
  # We want to retry, in case upstart-socket-bridge hasn't finished creating
  # our socket for us yet.
  n=0
  while [ $((n += 1)) -lt 10 ]; do
    output=$(chown root:${lp_group} "${CUPS_SOCKET}" 2>&1 && \
          chmod 770 "${CUPS_SOCKET}" 2>&1) && exit 0
    logger -t "${UPSTART_JOB}" "Failed to modify cups socket: ${output}"
    # Give a bit of delay before retrying.
    sleep 1
  done
  exit 1
end script
