# etc/init/cups-clear-state.conf
# Copyright 2016 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Perform CUPS session cleanup tasks."
author          "chromium-os-dev@chromium.org"

start on login-prompt-visible and started debugd

script
  dbus-send --system --print-reply --dest=org.chromium.debugd \
      /org/chromium/debugd org.chromium.debugd.CupsResetState || \
      logger -t "${UPSTART_JOB}" "Failed to reset cupsd"
end script
