# etc/init/swap.conf
# Copyright (c) 2012 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Setup compressed memory swap"
author          "chromium-os-dev@chromium.org"

env MIN_LOW_MEMORY_MARGIN=

start on starting system-services

exec /usr/share/cros/init/swap.sh start
