# etc/init/vm-nat.conf
# Copyright 2017 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Enable NAT for Termina VMs"
author        "chromium-os-dev@chromium.org"

# nat.conf lives in chromeos-base/chromeos-nat-init, and does the basic
# NAT setup.
start on started nat
stop on stopped nat

# This marks packets from _all_ interfaces starting with vmtap, since
# they all belong to termina, and will all want to be NAT'ed.
pre-start exec iptables -t mangle -A PREROUTING -i vmtap+ \
                        -j MARK --set-mark 1 -w

post-stop exec iptables -t mangle -D PREROUTING -i vmtap+ \
                        -j MARK --set-mark 1 -w
